# ikiwiki_supervisord.conf
[program:lighttpd]
command=/usr/sbin/lighttpd -D -f /etc/lighttpd/lighttpd.conf

[program:sshd]
command=/usr/sbin/sshd -D
